# bin/powersave.service
[Unit]
Description=Powersave
 
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/bin/powersave
 
[Install]
WantedBy=multi-user.target
